# openwrt-sdk-22.03.0-x86-generic_gcc-11.2.0_musl.Linux-x86_64/local-feeds/hamidrepo/udp_echo/scripts/etc/uci-defaults/99-udp-echo
#!/bin/sh
# uci-defaults scripts run once (on first boot after install) to set sane defaults.
# This script ensures /etc/config/udp_echo exists and enables the service.

# If config file missing, try copy from ROM overlay (where the package installed it)
[ -f /etc/config/udp_echo ] || cp /rom/etc/config/udp_echo /etc/config/udp_echo 2>/dev/null || true

# Enable the service so it starts on boot
/etc/init.d/udp-echo enable 2>/dev/null || true

exit 0
